# files/etc/rc.local
#!/bin/sh

bt-network -s nap br0 &
hciconfig hci0 up
hciconfig hci0 sspmode 1
hciconfig hci0 piscan
